# studio/micro-services/datasophon/datasophon-worker/src/main/resources/script/datasophon-worker
#!/bin/sh
# chkconfig: 35 95 1
usage="Usage: start.sh (start|stop|restart) <command> "

startStop=$1
installPath=$2

case $startStop in
  (start)
    sh ${installPath}/datasophon-worker/bin/datasophon-worker.sh start worker
    ;;
  (stop)
    sh ${installPath}/datasophon-worker/bin/datasophon-worker.sh stop worker
      ;;
  (status)
    sh ${installPath}/datasophon-worker/bin/datasophon-worker.sh status worker
      ;;
  (restart)
    sh ${installPath}/datasophon-worker/bin/datasophon-worker.sh restart worker
      ;;
  (*)
    echo $usage
    exit 1
    ;;

esac

echo "End $startStop $command."